# src/pwm-fan.service
[Unit]
Description=Configure PWM fan

[Service]
Type=oneshot
ExecStart=/usr/bin/rpi-temp-fan-pwm

[Install]
WantedBy=multi-user.target
